# etc/supervisord.d/apache.conf
[program:apache]
command=/usr/sbin/httpd -DFOREGROUND
