# etc/supervisord.d/mongo.conf
[program:mongo]
command=/opt/soltra/edge/bin/mongod -f /opt/soltra/edge/etc/mongod.conf
user=soltra-edge-mongod
